# gocat.socket
[Unit]
Description=Socket for Go HTTP Server Example

[Socket]
ListenStream=8080
